# doctor.service
[Unit]
Description=Doctor Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/doctor
ExecStart=/opt/doctor/venv/bin/python /opt/doctor/main.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target